# IntercomServer-RPi64.additionalFiles/etc/local.d/01-AddWlan0ToBr0.start
## We assume that wpa_supplicant is ready by now so we can add wlan0 to br0
ip link set dev wlan0 master br0